# pds/Professional_VOB/60-mercury/8khlinuxcab/mcucfg/mcu_8khm.conf
#for 8000H-M in redhat 6.0

start on startup
start on runlevel [2345]
stop  on runlevel [016]
stop  on shutdown

respawn

console output

exec /opt/mcu/start_mcu8khm.sh
