# pds/Professional_VOB/60-mercury/8kilinuxcab/mcucfg/mcu_8ki.conf
#for 8000i in ubuntu

start on startup
start on runlevel [2345]
stop  on runlevel [016]
stop  on shutdown

respawn

console output

exec /opt/mcu/ethdetect.sh
